# files/etc/uci-defaults/80-system-optimize
#!/bin/sh
# 系统优化配置

# 优化内存使用
uci -q batch <<-EOF >/dev/null
	set system.@system[0].cronloglevel='5'
	set system.@system[0].log_size='64'
	commit system
EOF

# 优化日志轮转
cat > /etc/logrotate.d/cups <<'EOF'
/var/log/cups/*.log {
    daily
    rotate 7
    compress
    delaycompress
    missingok
    notifempty
    create 644 root root
}
EOF

# 优化CUPS配置以提高性能
if [ -f /etc/cups/cupsd.conf ]; then
    # 添加性能优化设置
    cat >> /etc/cups/cupsd.conf <<'EOF'

# Performance optimizations
MaxJobs 50
MaxPrinterHistory 100
JobRetryLimit 5
JobRetryInterval 30
EOF
fi

# 设置系统时区为北京时间
uci set system.@system[0].timezone='CST-8'
uci set system.@system[0].zonename='Asia/Shanghai'
uci commit system

# 优化网络缓冲区
echo 'net.core.rmem_max = 262144' >> /etc/sysctl.conf
echo 'net.core.wmem_max = 262144' >> /etc/sysctl.conf
echo 'net.ipv4.tcp_rmem = 4096 87380 262144' >> /etc/sysctl.conf
echo 'net.ipv4.tcp_wmem = 4096 65536 262144' >> /etc/sysctl.conf

# 删除自己
rm -f /etc/uci-defaults/80-system-optimize
